# systemd/powersafe.service
[Unit]
Description=Powercut safe guard
Wants=powersafe.timer

[Service]
Type=oneshot
ExecStart=/sbin/powersafe

[Install]
RequiredBy=powersafe.timer
